# scripts/com.openclaw.virtual-office-site-monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.openclaw.virtual-office-site-monitor</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/zachli/.openclaw/workspace/virtual-office/scripts/site-monitor.sh</string>
    </array>
    <!-- Run every hour (3600 seconds) -->
    <key>StartInterval</key>
    <integer>3600</integer>
    <key>StandardOutPath</key>
    <string>/Users/zachli/.openclaw/workspace/virtual-office/scripts/site-monitor.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/zachli/.openclaw/workspace/virtual-office/scripts/site-monitor.log</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <false/>
    <key>LowPriorityIO</key>
    <true/>
</dict>
</plist>
